# services/sticker-api/deploy/majoin-stickers.service
[Unit]
Description=majoin sticker store API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/majoin-stickers
Environment=STICKER_ADMIN_KEY=CHANGE_ME_ADMIN_KEY
# uv resolves the venv from .venv/ created by `uv sync`.
ExecStart=/usr/local/bin/uv run uvicorn main:app --host 127.0.0.1 --port 8410
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
