# services/majoin/deploy/majoin.service
[Unit]
Description=majoin backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/majoin
# Secrets live in /opt/majoin/.env (not committed) — STICKER_ADMIN_KEY,
# DATABASE_URL. See deploy/.env.example for the format.
EnvironmentFile=/opt/majoin/.env
# uv resolves the venv from .venv/ created by `uv sync`.
ExecStart=/usr/local/bin/uv run uvicorn main:app --host 127.0.0.1 --port 8410
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
